# ext/root/anacrontab
# num@period  delay(unused) name script...
#
# uncomment to purge files older than 7days in TMP
#1@daily 0 purgetmp find /srv/data/tmp -type f -mtime +7 -delete > /dev/null
2@daily 0 update_bg php -d date.timezone='Europe/Paris' -f /srv/data/web/vhosts/physbook.fr/physbook/app/console update:bg --env=prod
@daily 0 notifications php -d date.timezone='Europe/Paris' -f /srv/data/web/vhosts/physbook.fr/physbook/app/console notifications --env=prod
@daily 0 brags_bucquage php -d date.timezone='Europe/Paris' -f /srv/data/web/vhosts/physbook.fr/physbook/app/console brags:bucquage --env=prod
@hourly 0 rezal_sync php -d date.timezone='Europe/Paris' -f /srv/data/web/vhosts/physbook.fr/physbook/app/console rezal:sync --env=prod
@daily 0 mysqldump_day mkdir -p /srv/data/tmp/backupSQL ; mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/daily_`date '+%F'`.databases.sql.gz ; rm -f /srv/data/tmp/backupSQL/daily_`date '+%F' --date '1 week ago'`.databases.sql.gz
@weekly 0 mysqldump_week mkdir -p /srv/data/tmp/backupSQL ; mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/weekly_`date '+%F'`.databases.sql.gz ; rm -f /srv/data/tmp/backupSQL/weekly_`date '+%F' --date '1 month ago'`.databases.sql.gz
@monthly 0 mysqldump_month mkdir -p /srv/data/tmp/backupSQL ; mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/monthly_`date '+%F'`.databases.sql.gz ; rm -f /srv/data/tmp/backupSQL/monthly_`date '+%F' --date '1 year ago'`.databases.sql.gz
@yearly 0 mysqldump_year mkdir -p /srv/data/tmp/backupSQL ; mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/yearly`date '+%F'`.databases.sql.gz
