# ext/root/anacrontab
# num@period  delay(unused) name script...
#
# uncomment to purge files older than 7days in TMP
#1@daily 0 purgetmp find /srv/data/tmp -type f -mtime +7 -delete > /dev/null
2@daily 0 update_bg php -f /srv/data/web/vhosts/physbook.fr/htdocs/app/console update:bg --env=prod
@daily 0 brags_bucquage php -f /srv/data/web/vhosts/physbook.fr/htdocs/app/console brags:bucquage --env=prod
@hourly 0 rezal_sync php -f /srv/data/web/vhosts/physbook.fr/htdocs/app/console rezal:sync --env=prod
echo $(date +`%Y-%m-%d %H:%i`) && " [START] ExportSQL"
@daily 0 purgedump_day find /srv/data/tmp/backupSQL -name "backupSQL_daily_*.sql.gz" -type f -ctime 1 -delete > /dev/null
@daily 0 mysqldump_day mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/backupSQL_daily_`date +%Y-%m-%d`.sql.gz
@weekly 0 purgedump_week find /srv/data/tmp/backupSQL -name "backupSQL_weekly_*.sql.gz" -type f -ctime 168 -delete > /dev/null
@weekly 0 mysqldump_week mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/backupSQL_weekly_`date +%Y-%m-%d`.sql.gz
@monthly 0 purgedump_month find /srv/data/tmp/backupSQL -name "backupSQL_monthly_*.sql.gz" -type f -ctime 5040 -delete > /dev/null
@monthly 0 mysqldump_month mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/backupSQL_monthly_`date +%Y-%m-%d`.sql.gz
@yearly 0 mysqldump_year mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/backupSQL_yearly_`date +%Y-%m-%d`.sql.gz
echo $(date +`%Y-%m-%d %H:%i`) && " [STOP] ExportSQL"