# root/anacrontab
# num@period  delay(unused) name script...
#
# uncomment to purge files older than 7days in TMP
#1@daily 0 purgetmp find /srv/data/tmp -type f -mtime +7 -delete > /dev/null
@daily 0 bucquage_brags php -f /srv/data/web/vhosts/physbook.fr/htdocs/app/console bucquage:brags --env=prod
4@daily 0 rezal_sync php -f /srv/data/web/vhosts/physbook.fr/htdocs/app/console rezal:sync
echo $(date +`%Y-%m-%d %H:%i`) && " [START] ExportSQL"
1@daily 0 purgedump_day find /srv/data/tmp/backupSQL -name "backupSQL_daily_*.sql.gz" -type f -ctime 1 -delete > /dev/null
1@daily 0 mysqldump_day mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/backupSQL_daily_`date +%Y-%m-%d`.sql.gz
2@weekly 0 purgedump_week find /srv/data/tmp/backupSQL -name "backupSQL_weekly_*.sql.gz" -type f -ctime 1 -delete > /dev/null
2@weekly 0 mysqldump_week mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/backupSQL_weekly_`date +%Y-%m-%d`.sql.gz
2@monthly 0 purgedump_month find /srv/data/tmp/backupSQL -name "backupSQL_monthly_*.sql.gz" -type f -ctime 1 -delete > /dev/null
2@monthly 0 mysqldump_month mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/backupSQL_monthly_`date +%Y-%m-%d`.sql.gz
3@yearly 0 mysqldump_year mysqldump -u dump -p'down' --opt physbook | /bin/gzip > /srv/data/tmp/backupSQL/backupSQL_yearly_`date +%Y-%m-%d`.sql.gz
echo $(date +`%Y-%m-%d %H:%i`) && " [STOP] ExportSQL"